# etc/zeppelin.conf
[program:zeppelin]
command=/opt/zeppelin/bin/zeppelin-daemon.sh start
